# container-files/etc/supervisor.d/squid.conf
[program:squid]
command = /usr/sbin/squid --foreground -f /etc/squid/squid.conf
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE